# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:java]
command=java -jar /app.jar

[program:nginx]
command=/etc/init.d/nginx start
